# cron/crontab
*/2 * 31 10 * test sjdklf sdkf    
1,3,5 15,18 */5 * * do_something
